# VenQore_Local/deploy/supervisor/venqore-horizon.conf
# ╔══════════════════════════════════════════════════════════════════╗
# ║  Supervisor Config — Laravel Horizon (VenQore Production)        ║
# ║  Phase 8 — Queue Workers for async job processing                ║
# ╚══════════════════════════════════════════════════════════════════╝
#
# Install:
#   sudo cp horizon.conf /etc/supervisor/conf.d/venqore-horizon.conf
#   sudo supervisorctl reread
#   sudo supervisorctl update
#   sudo supervisorctl start venqore-horizon:*
#
# Horizon dashboard: https://venqore.com/horizon
# Protect in routes/web.php with Gate::define('viewHorizon', ...) or auth middleware.

[program:venqore-horizon]
process_name=%(program_name)s
command=php /var/www/venqore/artisan horizon
directory=/var/www/venqore
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/venqore-horizon.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600       ; Wait up to 1 hour before SIGKILL on graceful stop
killasgroup=true
stopasgroup=true

# ── Horizon Queues (defined in config/horizon.php) ────────────────────────────
# These are the Horizon workers that match the queues already used in the app:
#
#   default     → General jobs (email, notifications, growth engine scans)
#   provision   → New tenant DB seeding (ProvisionTenantJob — fast, short)
#   webhooks    → Lemon Squeezy + AppSumo webhook processing (reliable)
#   reports     → Heavy async report generation (slow, background)
#
# Horizon auto-balances workers based on queue depth.
# See: config/horizon.php for full worker configuration.

# ── Scheduler (cron instead of Supervisor) ────────────────────────────────────
# Add this to www-data's crontab: sudo crontab -u www-data -e
#
#   * * * * * cd /var/www/venqore && php artisan schedule:run >> /dev/null 2>&1
#
# This runs all scheduled commands from routes/console.php including:
#   - demo:reset          (daily 04:00)
#   - tenants:send-trial-reminders   (daily 09:00)
#   - tenants:process-expired-trials (hourly)
#   - staff:generate-daily-summaries (daily 00:05)
#   - growth:analyze                 (daily 09:00)
